# debian/preinst
#!/bin/sh
# preinst script for logstash-forwarder
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <new-preinst> `install'
#        * <new-preinst> `install' <old-version>
#        * <new-preinst> `upgrade' <old-version>
#        * <old-preinst> `abort-upgrade' <new-version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    upgrade)
    # Manually stop logstash-forwarder when upgrading
    # from an upstream package which has the binary
    # reside in /opt
    if dpkg --compare-versions 0.3.1 eq 0.3.1; then
        invoke-rc.d logstash-forwarder stop || true
    fi
    ;;
    install)
    # do nothing
    ;;
    abort-upgrade)
    ;;
    *)
        echo "preinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
